# supervisord.conf

[program:apache_profiles]
command=apache2 -DFOREGROUND
stdout_logfile=syslog
stderr_logfile=syslog

